# target/linux/wm2500rp/base-files/etc/init.d/network
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=40
STOP=40

start() {
	echo "Disable script at present --> network start" && exit 0
}

stop() {
	echo "Disable script at present --> network stop"
}

restart() {
	stop
	start
}
