# target/linux/x6200/base-files/etc/init.d/detbutton
#!/bin/sh /etc/rc.common
START=15

start() {
	# Start detbutton
	gpio &
}

stop() {
	killall gpio
}

restart() {
	stop
	start
}
