# target/linux/x6200/base-files/etc/init.d/net-lan
#!/bin/sh /etc/rc.common
#START=30

DEFAULT_IP="192.168.1.250"
LLD2D="lld2d"

print_dhcpd_conf(){
    cat <<EOF
pidfile /var/run/udhcpd.pid
start $2
end $3
interface $1
remaining yes
auto_time 5
lease_file /tmp/udhcpd.leases
option subnet $4
option router $5
option dns $5
option lease $6
EOF
}
 
pull_down_eth() {
	for i in 1 2 3 4 ; do
		old_val=`et-tool phyrd $i 0x1e`
		let gt3bit=old_val/16
		let lt3bit=old_val%16%8
		let new_val=gt3bit*16+8+lt3bit
		et-tool phywr $i 0x1e $new_val
	done
}
pull_up_eth() {
	for i in 1 2 3 4 ; do
		old_val=`et-tool phyrd $i 0x1e`
		let gt3bit=old_val/16
		let lt3bit=old_val%16%8
		let new_val=gt3bit*16+lt3bit
		et-tool phywr $i 0x1e $new_val
	done
}

start_dhcpd() {
    [ ! -f /tmp/udhcpd.leases ] && touch /tmp/udhcpd.leases
    print_dhcpd_conf "$BR_IF" "$($CONFIG get dhcp_start)" "$($CONFIG get dhcp_end)" "$($CONFIG get lan_netmask)" "$($CONFIG get lan_ipaddr)" "$($CONFIG get lan_lease)" > /tmp/udhcpd.conf

    count=1
    while :
    do
	lease="$($CONFIG get reservation$count)"
	if [ "x$lease" = "x" ]; then
	    break;
	fi
	count=`expr $count + 1`
	echo "static_lease $lease" >> /tmp/udhcpd.conf
    done

    udhcpd /tmp/udhcpd.conf

    # Fix bug 34562: [LOG]There is no log for the router assigns a DHCP to a LAN client while DUT boot up.
    # Should restart the ethernet phy.

    # Pull down the Ethernet PHY
    pull_down_eth
    sleep 6
    # Up the Ethernet PHY
    pull_up_eth
}

start() {

	if [ "x`$CONFIG get lan_dhcp`" = "x1" ]; then
		# whatever run udhcpc, then udhcpd depends on if we get an IP
		udhcpc -b -i $BR_IF -r $DEFAULT_IP -h "$($CONFIG get netbiosname)"	
	else
		# set static ip
		lan_ipaddr=$($CONFIG get lan_ipaddr)
		lan_gateway=$($CONFIG get lan_gateway)
		lan_netmask=$($CONFIG get lan_netmask)
		lan_broadcast=$($CONFIG get lan_broadcast)
		[ -n "$lan_broadcast" ] && BROADCAST="broadcast $lan_broadcast"
		[ -n "$lan_netmask" ] && NETMASK="netmask $lan_netmask"
		ifconfig br0 $lan_ipaddr $BROADCAST $NETMASK

		if [ -n "$lan_gateway" ]; then
			route del default gw 0.0.0.0 dev br0 > /dev/null 2>&1
			route add default gw $lan_gateway dev br0
		fi

		# turn off hijack all
		echo 0 > /proc/sys/net/ipv4/dns_hijack
	fi


	# Attach Device
	mkdir -p /tmp/netscan
	/usr/sbin/net-scan

#	# Restart UPnP
#	[ ! -f /tmp/boot_status ] && /sbin/cmdupnp restart

	# Start Router Debugging Mode ...
#	. /lib/network/RtDebug.sh
#	start_RtDebugMode

#	# Apply static route    
#	/sbin/cmdroute start

	# Start lld2d
#	echo "$($CONFIG get netbiosname)" > /proc/sys/kernel/hostname
	$LLD2D $BR_IF
	
#	sleep 2
	#For ipv6 ready log test , fixed ip setting do not enable radvd and dhcpv6s
#	local wan6_type=$($CONFIG get ipv6_type)
#	local logo_test=$($CONFIG get endis_ipv6_logo_test)
#	if [ "x$logo_test" = "x1" -a "x$wan6_type" = "xfixed" ];then
#		/etc/net6conf/6fixed lan
#	fi

}

stop() {
	killall udhcpc
	killall udhcpd
	killall net-scan
	killall lld2d

	# PHY link will be pulled low some seconds to force transition to reboot state 
	# and generating DHCP request and Discovery protocol and address refresh in the 
	# devices connected to the NETGEAR Local Area Network ports.
	#switch_pwr 6
}

restart() {
	stop
	start
}
