# target/linux/x6200/base-files/etc/init.d/wireless
#!/bin/sh /etc/rc.common
START=25
boot() {
	/sbin/rpt-wlan start
	#(sleep 3 && wlconf eth2 start)&
	# generate random ssid and psk
	/sbin/genr_ssid_key
}
