# package/base-files/realtek-base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
ttyS0::respawn:/bin/ash --login
